# infra/teses.socket
[Unit]
Description=teses socket

[Socket]
ListenStream=/run/teses.sock

[Install]
WantedBy=sockets.target